# files/systemd/kiosk-api.service
[Unit]
Description=Kiosk FastAPI
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/kiosk
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/kiosk/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8989
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target